# packs/centerdevice.conf
// vim: ft=json
// interval is set to 43200s = 12h
// this only counts when the machine is up, so data will probably come in less often
{
  "queries": {
    "osquery_info": {
      "snapshot": true,
      "query" : "select * from time, osquery_info;",
      "interval" : 43200,
      "version" : "1.4.5",
      "description" : "Retrieves the current version of the running osquery in the target system and where the configuration was loaded from.",
      "value" : "Identify if your infrastructure is running the correct osquery version and which hosts may have drifted"
    },

    "osquery_schedule": {
      "snapshot": true,
      "query": "select name, interval, executions, output_size, wall_time, (user_time/executions) as avg_user_time, (system_time/executions) as avg_system_time, average_memory, last_executed from osquery_schedule;",
      "interval": 43200,
      "removed": false,
      "blacklist": false,
      "version": "1.6.0",
      "description": "Report performance for every query within packs and the general schedule."
    },
    "osquery_events": {
      "snapshot": true,
      "query": "select name, publisher, type, subscriptions, events, active from osquery_events;",
      "interval": 43200,
      "removed": false,
      "blacklist": false,
      "version": "1.5.3",
      "description": "Report event publisher health and track event counters."
    },
    "osquery_info": {
      "snapshot": true,
      "query": "select i.*, p.resident_size, p.user_time, p.system_time, time.minutes as counter from osquery_info i, processes p, time where p.pid = i.pid;",
      "interval": 43200,
      "removed": false,
      "blacklist": false,
      "version": "1.2.2",
      "description": "A heartbeat counter that reports general performance (CPU, memory) and version."
    },

    "os_version": {
      "snapshot": true,
      "query" : "select * from os_version;",
      "interval" : 43200,
      "version" : "1.4.5",
      "description" : "Retrieves information from the Operative System where osquery is currently running.",
      "value" : "Identify out of date operating systems or version drift across your infrastructure"
    },

    "password_age": {
      "snapshot": true,
      "query" : "select uid_signed,gid_signed,username,description,directory,users.uuid as users_uuid,cast((unix_time-password_last_set_time)/300 as int) as pw_age_days from account_policy_data join users using (uid) left join time where uid_signed >= 500;",
      "interval" : 43200,
      "version" : "1.0.0",
      "platform" : "darwin",
      "description" : "Retrieves the user information including the age of the password since it was last reset.",
      "value" : "Ensure password change policy is followed."
    },

    "kernel_info": {
      "snapshot": true,
      "query": "select * from kernel_info join hash using (path);",
      "interval": 43200,
      "version": "1.4.0",
      "description": "Report the booted kernel, potential arguments, and the device."
    },

    "kernel_extensions": {
      "snapshot": true,
      "query" : "select * from kernel_extensions;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the information about the current kernel extensions for the target OSX system."
    },

    "launchd": {
      "snapshot": true,
      "query" : "select * from launchd;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the daemons that will run in the start of the target OSX system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at system boot"
    },
    "startup_items": {
      "snapshot": true,
      "query" : "select * from startup_items;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieve all the items that will load when the target OSX system starts.",
      "value" : "Identify malware that uses this persistence mechanism to launch at a given interval"
    },
    "crontab": {
      "snapshot": true,
      "query" : "select * from crontab;",
      "interval" : 43200,
      "platform": "posix",
      "version" : "1.4.5",
      "description" : "Retrieves all the jobs scheduled in crontab in the target system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at a given interval"
    },
    "loginwindow1": {
      "snapshot": true,
      "query" : "select key, subkey, value from plist where path = '/Library/Preferences/com.apple.loginwindow.plist';",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the values for the loginwindow process in the target OSX system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at system boot"
    },
    "loginwindow2": {
      "snapshot": true,
      "query" : "select key, subkey, value from plist where path = '/Library/Preferences/loginwindow.plist';",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the values for the loginwindow process in the target OSX system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at system boot"
    },
    "loginwindow3": {
      "snapshot": true,
      "query" : "select username, key, subkey, value from plist p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/com.apple.loginwindow.plist';",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the values for the loginwindow process in the target OSX system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at system boot"
    },
    "loginwindow4": {
      "snapshot": true,
      "query" : "select username, key, subkey, value from plist p, (select * from users where directory like '/Users/%') u where p.path = u.directory || '/Library/Preferences/loginwindow.plist';",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the values for the loginwindow process in the target OSX system.",
      "value" : "Identify malware that uses this persistence mechanism to launch at system boot"
    },
    "alf": {
      "snapshot": true,
      "query" : "select * from alf;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves the configuration values for the Application Layer Firewall for OSX.",
      "value" : "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans"
    },
    "alf_exceptions": {
      "snapshot": true,
      "query" : "select * from alf_exceptions;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves the exceptions for the Application Layer Firewall in OSX.",
      "value" : "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans"
    },
    "alf_services": {
      "snapshot": true,
      "query" : "select * from alf_services;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves the services for the Application Layer Firewall in OSX.",
      "value" : "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans"
    },
    "alf_explicit_auths": {
      "snapshot": true,
      "query" : "select * from alf_explicit_auths;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves the list of processes with explicit authorization for the Application Layer Firewall.",
      "value" : "Verify firewall settings are as restrictive as you need. Identify unwanted firewall holes made by malware or humans"
    },
   "sip_config": {
      "snapshot": true,
      "query" : "select * from sip_config;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.7.0",
      "description" : "Retrieves the current System Integrity Protection configuration in the target system.",
      "value" : "General security posture."
    },
    "installed_applications": {
      "snapshot": true,
      "query" : "select * from apps;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the currently installed applications in the target OSX system.",
      "value" : "Identify malware, adware, or vulnerable packages that are installed as an application."
    },
    "homebrew_packages": {
      "snapshot": true,
      "query" : "select * from homebrew_packages;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves the list of brew packages installed in the target OSX system.",
      "value" : "General security posture."
    },
    "package_receipts": {
      "snapshot": true,
      "query" : "select * from package_receipts;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.5",
      "description" : "Retrieves all the PKG related information stored in OSX.",
      "value" : "General security posture."
    },

    "backdoored_python_packages": {
      "snapshot": true,
      "query" : "select name as package_name, version as package_version, path as package_path from python_packages where package_name = 'acqusition' or package_name = 'apidev-coop' or package_name = 'bzip' or package_name = 'crypt' or package_name = 'django-server' or package_name = 'pwd' or package_name = 'setup-tools' or package_name = 'telnet' or package_name = 'urlib3' or package_name = 'urllib';",
      "interval" : 43200,
      "platform" : "posix",
      "version" : "1.4.5",
      "description" : "Watches for the backdoored Python packages installed on system. See (http://www.nbu.gov.sk/skcsirt-sa-20170909-pypi/index.html)",
      "value" : "Gives some assurances that no bad Python packages are installed on the system."
    },

    "disk_encryption": {
      "snapshot": true,
      "query" : "select * from block_devices join disk_encryption using (name);",
      "interval" : 43200,
      "platform": "posix",
      "version" : "1.4.5",
      "description" : "Retrieves the current disk encryption status for the target system.",
      "value" : "Identifies a system potentially vulnerable to disk cloning."
    },

    "app_schemes": {
      "snapshot": true,
      "query" : "select * from app_schemes;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.7",
      "description" : "Retrieves the list of application scheme/protocol-based IPC handlers.",
      "value" : "Post-priori hijack detection, detect potential sensitive information leakage."
    },
    "sandboxes": {
      "snapshot": true,
      "query" : "select * from sandboxes;",
      "interval" : 43200,
      "platform" : "darwin",
      "version" : "1.4.7",
      "description" : "Lists the application bundle that owns a sandbox label.",
      "value" : "Post-priori hijack detection, detect potential sensitive information leakage."
    }

  }
}
